# services/prometheus-pp.service
[Unit]
Description=Prometheus BitMEX Grid Trading Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/prometheus
ExecStart=/home/ubuntu/prometheus/start-pp.sh
KillMode=process
Restart=on-failure
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=prometheus-pp
LimitNOFILE=4096
IgnoreSIGPIPE=false
Environment=NODE_ENV=production
Environment=BITMEX_API_KEY=
Environment=BITMEX_API_SECRET=
Environment=TRADING_SYMBOL=XBTUSD
Environment=INFLUX_HOST=
Environment=INFLUX_TOKEN=
Environment=INFLUX_DATABASE=prometheus
Environment=INFLUX_ENABLED=true

[Install]
WantedBy=multi-user.target 